# pkg/debian/disco-master.prerm
#!/bin/sh

set -e

invoke-rc.d disco-master stop
