# pkg/debian/disco-master.preinst
#!/bin/sh

set -e

[ -e /etc/init.d/disco-master ] && invoke-rc.d disco-master stop
