# pkg/debian/disco-node.preinst
#!/bin/sh

[ -e /etc/init.d/disco-node ] && invoke-rc.d disco-node stop
exit 0
