# debian/disco-master.preinst
#!/bin/sh

invoke-rc.d disco-node stop
invoke-rc.d disco-master stop
