# debian/disco-master.prerm
#!/bin/sh

invoke-rc.d disco-master stop
