# contrib/docker/supervisor.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:disco]
command=/usr/bin/disco nodaemon -v
user=disco
environment=HOME="/home/disco"

